# crontab
01 00 * * * bash /scripts/log_analyzer.sh
